# raspberry-pi/signreader-capture.service
[Unit]
Description=SignReader Pi Webcam Capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/signreader-pi
EnvironmentFile=/opt/signreader-pi/.env
ExecStart=/opt/signreader-pi/venv/bin/python capture.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
